# DEBIAN/postinst
#!/bin/bash
systemctl enable getty@tty1.service
exit 0