# DEBIAN/postinst
#!/bin/bash
set -e

systemctl enable getty@tty1.service
exit 0
